# deploy/jonglerie.service
[Unit]
Description=Jonglerie — service web d'analyse de jongle (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/jonglerie
Environment=JONGLE_DATA=/opt/jonglerie/data
Environment=JONGLE_MODEL=/opt/jonglerie/models/yolov8m.onnx
ExecStart=/opt/jonglerie/.venv/bin/uvicorn app:app --app-dir /opt/jonglerie/server --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
